# crontab
@reboot /bin/bash /home/pi/MQTT_Thermostat/tmux.sh
